# configs/supervisor.conf
[program:worker]
command=python /var/www/derpi.nl/manage.py runworker
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log
directory=/var/www/derpi.nl
[program:daphne]
command=daphne derpi.asgi:channel_layer 
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/daphne.err.log
stdout_logfile=/var/log/daphne.out.log
directory=/var/www/derpi.nl
